# rootfs/etc/s6/.s6-svscan/finish
#!/usr/bin/env bash
for file in /etc/s6/*/finish; do
   $file
done